# contrib/tlstun.initd
#!/sbin/openrc-run
# Copyright 1999-2021 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

name="${RC_SVCNAME}"
command="/usr/bin/tlstun"
pidfile="/run/${RC_SVCNAME}.pid"
command_args="${TLSTUN_ARGS}"
command_user="nobody"
command_background="yes"
